# smf/manifests/nginx.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- Copyright (c) 2012, Joyent, Inc. All rights reserved. -->
<service_bundle type="manifest" name="nginx" >
    <service name="sds/application/mako" type="service" version="1" >
	<single_instance />
        <create_default_instance enabled="true" />
	<single_instance />

	<dependency name="network"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/network/physical" />
	</dependency>

	<dependency name="filesystem"
		    grouping="require_all"
		    restart_on="error"
		    type="service">
	    <service_fmri value="svc:/system/filesystem/local" />
	</dependency>

        <exec_method type="method" name="start" exec="/opt/smartdc/mako/nginx/sbin/nginx -p /opt/smartdc/mako/nginx/ &amp;" timeout_seconds="11" />

        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" />

        <template >
            <common_name >
                <loctext xml:lang="C" >Joyent nginx HTTP server</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
